# allsite-hcm.service
[Unit]
Description=Allsite HCM DDoS Protection Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/allsite-hcm
ExecStart=/usr/bin/python3 /opt/allsite-hcm/Allsite-HCM_Service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
